# livecd/var/lib/pacman/local/mdadm-3.2.3-1/files
%FILES%
etc/
etc/mdadm.conf
etc/rc.d/
etc/rc.d/mdadm
lib/
lib/initcpio/
lib/initcpio/hooks/
lib/initcpio/hooks/mdadm
lib/initcpio/hooks/raid
lib/initcpio/install/
lib/initcpio/install/mdadm
lib/initcpio/install/mdadm_udev
lib/udev/
lib/udev/rules.d/
lib/udev/rules.d/64-md-raid.rules
sbin/
sbin/mdadm
sbin/mdassemble
sbin/mdmon
usr/
usr/share/
usr/share/man/
usr/share/man/man4/
usr/share/man/man4/md.4.gz
usr/share/man/man5/
usr/share/man/man5/mdadm.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/mdadm.8.gz
usr/share/man/man8/mdmon.8.gz

%BACKUP%
etc/mdadm.conf	00cbed931db4f15b6ce49e3e7d433966

