# livecd/var/lib/pacman/local/rsync-3.0.9-2/files
%FILES%
etc/
etc/rc.d/
etc/rc.d/rsyncd
etc/rsyncd.conf
etc/xinetd.d/
etc/xinetd.d/rsync
usr/
usr/bin/
usr/bin/rsync
usr/lib/
usr/lib/rsync/
usr/lib/rsync/rrsync
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/rsync.1.gz
usr/share/man/man5/
usr/share/man/man5/rsyncd.conf.5.gz

%BACKUP%
etc/rsyncd.conf	4395c0591638349b1a7aeaaa4da1f03a
etc/xinetd.d/rsync	ea3e9277dc908bc51f9eddc0f6b935c1

